# user-settings/uci-defaults/lan-ip-6.uci.default
#!/bin/sh

uci -q get network && {
	uci -q batch <<-EOF >/dev/null
		set network.lan.ipaddr='192.168.6.1'
		commit network
	EOF
}

exit 0
