# user-settings/files/etc/uci-defaults/15_tailscale
#!/bin/sh

uci -q get tailscale && {
    [ -x "/etc/init.d/tailscale" ] && /etc/init.d/tailscale disable
}

exit 0
