# user-settings/files/etc/uci-defaults/15_upnpd
#!/bin/sh

uci -q get upnpd && {
	uci -q batch <<-EOF >/dev/null
		set upnpd.config.enabled='1'
		commit upnpd
	EOF
}

exit 0
